# docker/etc/services.d/maltrail-sensor/run
#!/command/with-contenv bash
set -e

echo "starting maltrail sensor.."

cd ${VENV}
python sensor.py