# docker/etc/services.d/maltrail-server/run
#!/command/with-contenv bash
set -e

echo 'starting maltrail webserver..'

cd ${VENV}
python server.py

